# deploy/valora-supervisor.conf
[program:valora-api]
command=/usr/local/bin/python -m uvicorn valora_backend.main:app --host 127.0.0.1 --port 8003
directory=/app
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:valora-web]
command=/bin/sh -c "sleep 2 && exec /usr/local/bin/node ./node_modules/next/dist/bin/next start -H 0.0.0.0 -p \"${PORT:-3000}\""
directory=/app/frontend
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
